# unison/docker/node2/conf/serf.conf
[program:serf]
command=/usr/bin/serf agent -config-dir=/etc/serf 
numprocs=1
autostart=true
autorestart=true
