# unison/docker/node2/conf/unison.conf
[program:unison]
command=/usr/bin/unison -socket 7654
stdout_events_enabled=true
stderr_evetns_enabled=true

