# unison/docker/central/conf/dnsmasq.conf
[program:dnsmasq]
command=/usr/sbin/dnsmasq -d -C /etc/dnsmasq.conf 
numprocs=1
autostart=true
autorestart=true
