# unison/docker/central/conf/serf_add_node.conf
[program:serf-add-node]
command=/add-node.sh
autostart=false 
autorestart=false 
startretries=0 
redirect_stderr=false 
redirect_stderr=false 
