# unison/docker/central/conf/tinyweb.conf
[program:tinyweb]
command=python /conf/droopy.py -d /srv/data --dl  
stdout_events_enabled=true
stderr_evetns_enabled=true

